# conf/supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port = 9001
#username = user # Basic auth username
#password = pass # Basic auth password

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=/usr/sbin/nginx "-g" "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
